# library-service/supervisord.conf
[supervisord]
nodaemon=true

[program:springboot]
directory=/root/
command=python3 /root/server.py

[program:tor]
command=service tor start